# conf/crontab
# предусматриваем перезагрузку сервера
@reboot $USERDIR/app start
